# systemd/odroid/bush-t2v.service
[Unit]
Description=text-to-verse with MQTT
After=chromadb.service ollama.service
Requires=chromadb.service ollama.service

[Service]
Type=simple
User=odroid
Environment=T2V_BIN=/home/odroid/repos/t2v/target/release/text-to-verse
Environment=AFFECTS_DIR=/home/odroid/repos/t2v/templates/affects
ExecStart=/usr/bin/python3 /home/odroid/repos/bushglue/t2v-service.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
