# systemd/odroid/bush-tts.service
[Unit]
Description=Text-to-Speech verse reader
After=network-online.target user@1000.service

[Service]
Type=simple
User=odroid
ExecStart=/home/odroid/bushglue/.venv/bin/bush-tts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
