# systemd/odroid/bush-t2v.service
[Unit]
Description=text-to-verse with MQTT
After=chromadb.service ollama.service
Requires=chromadb.service ollama.service

[Service]
Type=simple
User=odroid
Environment=T2V_BIN=/home/odroid/bushglue/t2v/target/release/text-to-verse
Environment=AFFECTS_DIR=/home/odroid/bushglue/t2v/templates/affects
ExecStart=/home/odroid/bushglue/.venv/bin/bush-t2v
MemoryMax=1500M
MemoryHigh=1200M
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
