# systemd/odroid/bush-standby.service
[Unit]
Description=Bush standby mode (loop audio, fire control web UI)
After=network-online.target mosquitto.service sound.target

[Service]
Type=simple
ExecStart=/home/odroid/bushglue/utils/bush-standby
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
