# systemd/odroid/bush-audio-agent.service
[Unit]
Description=Audio device discovery agent
After=network-online.target user@1000.service

[Service]
Type=simple
User=odroid
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=PULSE_SERVER=unix:/run/user/%U/pulse/native
ExecStart=/usr/bin/python3 %h/repos/bushglue/audio-agent.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
