# systemd/odroid/chromadb.service
[Unit]
Description=ChromaDB Vector Database
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=odroid
ExecStart=%h/repos/t2v_chroma/.venv/bin/chroma run --path %h/repos/t2v_chroma/chromadb
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
