# systemd/bush-stt.service
[Unit]
Description=Speech-to-Text MQTT Publisher
After=network-online.target

[Service]
Type=simple
User=ubuntu
Environment=PULSE_SERVER=unix:/mnt/wslg/PulseServer
Environment=STT_DIR=%h/repos/speech-to-text
Environment=STT_MODEL=%h/repos/speech-to-text/models/en-us
ExecStart=/usr/bin/python3 %h/repos/bushglue/stt-service.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
