# systemd/bush-sentiment.service
[Unit]
Description=Sentiment Analysis with MQTT Fire Control
After=network.target

[Service]
Type=simple
User=ubuntu
ExecStart=/home/ubuntu/bbsentimentqq-venv/bin/python3 /mnt/c/Users/EB/bushglue/scripts/sentiment-service.py
WorkingDirectory=/mnt/c/Users/EB/bushglue/scripts
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
