# systemd/bush-t2v.service
[Unit]
Description=text-to-verse with MQTT
After=chromadb.service ollama.service
Requires=chromadb.service ollama.service

[Service]
Type=simple
User=ubuntu
ExecStart=/usr/bin/python3 /mnt/c/Users/EB/bushglue/scripts/t2v-service.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
